# roles/systemd/files/systemd/user/mitmweb.service
[Unit]
Description=Local HTTP proxy
PartOf=graphical-session.target

[Service]
ExecStart=mitmweb -p 10000 --web-port 10001
Restart=on-failure
