# roles/systemd/files/systemd/user/polybar.service
[Unit]
Description=Polybar Panel
PartOf=graphical-session.target

[Service]
ExecStart=polybar --reload main
Restart=on-failure
