# roles/systemd/files/systemd/user/polybar.service
[Unit]
Description=Polybar Panel
PartOf=graphical-session.target

[Service]
ExecStart=%h/.config/polybar/launch.sh
ExecReload=%h/.config/polybar/launch.sh
ExecStop=killall --quiet --user polybar
Restart=on-failure
