# roles/systemd/files/systemd/user/restic-prune.timer
[Unit]
Description=Restic prune old snapshots on calendar

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
