# roles/systemd/files/systemd/user/sxhkd.service
[Unit]
Description=Simple X Hotkey Daemon (sxhkd)
PartOf=graphical-session.target

[Service]
ExecStart=sxhdk
ExecReload=kill -USR1 $MAINPID
Restart=on-failure
