# roles/systemd-user/files/systemd/user/mpd-notification.service.d/override.conf
[Unit]
After=
After=mpd.service network.target network-online.target
After=dunst.service
Wants=dunst.service
Wants=mpd.service

[Service]
ExecStart=
ExecStart=mpd-notification -m "/home/Music" -t 8
