# roles/systemd-user/files/systemd/user/mpc-update.service
[Unit]
Description=Update mpd music library
PartOf=graphical-session.target
Wants=mpd.service
After=mpd.service

[Service]
Type=oneshot
ExecStart=mpc update
