# roles/systemd-user/files/systemd/user/picom.service
[Unit]
Description=picom compositor
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/picom
Slice=session.slice

Restart=always
