# roles/systemd-user/files/systemd/user/restic.timer
[Unit]
Description=Run Restic on calendar

[Timer]
OnCalendar=Mon,Tue,Wed,Thu,Fri *-*-* 17:45:00

[Install]
WantedBy=timers.target
