# roles/systemd-user/files/systemd/user/xrdb-load.service
[Unit]
Description=Load Xresources
PartOf=graphical-session.target

[Service]
Type=oneshot
ExecStart=xrdb -load "$XDG_CONFIG_HOME/xorg/.Xresources"
