# roles/systemd-user/files/systemd/user/xsettingsd.service
[Unit]
Description=XSETTINGS-protocol daemon
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/xsettingsd
Slice=session.slice
